# docs/llm-backend-plan/examples/vel-openai-oauth.service
[Unit]
Description=Vel openai-oauth local proxy
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/src/vel
ExecStart=%h/src/vel/services/llm-gateway/scripts/run-openai-oauth.sh
Restart=on-failure
RestartSec=2
Environment=HOST=127.0.0.1
Environment=PORT=8014

[Install]
WantedBy=default.target
